# zeroclaw/host/saviaclaw-headless.service
[Unit]
Description=SaviaClaw Headless Agent — autonomous PM agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=monica
WorkingDirectory=/home/monica/claude
ExecStart=/usr/bin/python3 -m zeroclaw.host.saviaclaw_headless
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/monica/.opencode/bin:/home/monica/.local/bin:/usr/local/bin:/usr/bin:/bin

NoNewPrivileges=true
PrivateTmp=true
MemoryMax=512M

[Install]
WantedBy=multi-user.target
